# docker/noVNC/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 8080 starcraft:5901
autorestart=true
